# Usabilidade/useskill-capture/config/nodejs-init.conf
start on startup
exec forever start /home/useskill/useskill_capture/app-fm.js
